# intranet.service
[Unit]
Description=Intranet ES-SERVIDOR Platform
After=network.target

[Service]
Type=simple
User=alfreire
Group=alfreire
WorkingDirectory=/opt/intranet
EnvironmentFile=-/opt/intranet/.env
Environment="PATH=/opt/intranet/venv/bin"
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/intranet/venv/bin/gunicorn --worker-class eventlet --workers 1 --bind 127.0.0.1:5000 --access-logfile /opt/intranet/logs/access.log --error-logfile /opt/intranet/logs/error.log wsgi:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
