# intranet.service
[Unit]
Description=Intranet ES-SERVIDOR Platform
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/intranet
EnvironmentFile=-/opt/intranet/.env
Environment="PATH=/opt/intranet/venv/bin"
ExecStart=/opt/intranet/run.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
